# app/cron/crontab
* * * * * root python3 /opt/scripts/overpass_data_update.py >> /var/log/cron.log
# Empty line